# _config/autostart/setup_displays.desktop
[Desktop Entry]
Type=Application
Name=Setup Displays
Exec=setup-displays
X-GNOME-Autostart-enabled=true
Comment=Display managment
